# com.teachce.crm-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.teachce.crm-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/timothysepulvado/Teach/CRM_integration/sync.sh</string>
        <string>full</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/timothysepulvado/Teach/CRM_integration</string>

    <!-- Run daily at 6:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/timothysepulvado/Teach/CRM_integration/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/timothysepulvado/Teach/CRM_integration/logs/launchd_error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
